# Config/LobstaListen.conf
[program:LobstaListen]
directory = /home/ubuntu/LobstaListen
umask = 002
command = /home/ubuntu/LobstaListen/LobstaListen.py
stopwaitsecs = 30
priority = 100

